# max-bridge.service
[Unit]
Description=MAX ↔ OpenClaw Bridge
After=network.target

[Service]
Type=simple
User=aurum
WorkingDirectory=/opt/max-bridge
EnvironmentFile=/opt/max-bridge/.env
ExecStart=/usr/bin/python3 -u /opt/max-bridge/max-bridge.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
